# ndi-receiver.service
[Unit]
Description=NDI Receiver LED Screen Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=catatumbo
WorkingDirectory=/home/catatumbo/led_test
Environment="PYTHONUNBUFFERED=1"
Environment="XDG_RUNTIME_DIR=/run/user/1000"
ExecStart=/usr/bin/python3 ndi_receiver.py --config config.led_screen.json --bridge-url ws://localhost:8081
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

